# backend/web_app/build/wsgi_setup/agora-http-api.service
[Unit]
Description=agora rest api
After=network.target

[Service]
User=cloud-user
WorkingDirectory=/home/cloud-user/plateform/agora/backend/web_app
ExecStart=/usr/local/bin/gunicorn --preload wsgi_http:app
Restart=always


[Install]
WantedBy=multi-user.target
